# ~/Library/LaunchAgents/com.tk-auto-dl.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.tk-auto-dl.daily</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/walhalax/Library/Mobile Documents/com~apple~CloudDocs/Coding/VSCode/Project/tk-auto-dl/src/cli_runner.py</string>
    </array>
    
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>2</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>
    
    <key>StandardOutPath</key>
    <string>/Users/walhalax/Library/Mobile Documents/com~apple~CloudDocs/Coding/VSCode/Project/tk-auto-dl/logs/daily.stdout.log</string>
    
    <key>StandardErrorPath</key>
    <string>/Users/walhalax/Library/Mobile Documents/com~apple~CloudDocs/Coding/VSCode/Project/tk-auto-dl/logs/daily.stderr.log</string>
    
    <key>RunAtLoad</key>
    <false/>
    
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>